# roles/pihole/templates/dns.service.j2
[Unit]
Description=Dockerized DNS service for Raspberry Pi
After=docker.service
Requires=docker.service

[Service]
WorkingDirectory=/opt/dns
ExecStart=docker-compose -f /opt/dns/docker-compose.yml up -d
ExecStop=docker-compose -f /opt/dns/docker-compose.yml down

[Install]
WantedBy=multi-user.target