# roles/pihole/templates/dns.service.j2
[Unit]
Description=Dockerized DNS service for Raspberry Pi
After=docker.service
Requires=docker.service

[Service]
WorkingDirectory=/opt/dns
ExecStart=/usr/bin/docker-compose up --force-recreate
ExecReload=/usr/bin/docker-compose restart
ExecStop=/usr/bin/docker-compose stop

[Install]
WantedBy=multi-user.target